# roles/physical/templates/udev/rules.d/90-monitor-hotplug.rules.j2
ACTION=="change", KERNEL=="card0", SUBSYSTEM=="drm", RUN+="{{ ansible_user_dir }}/.local/bin/monitor_hotplug.sh"
